# errbit.cron
0  0  *  *  * root cd /home/app/errbit && RAILS_ENV=production bundle exec rake errbit:db:clear_resolved
